# systemd/gcp-watchdog-web.service
[Unit]
Description=GCP Watchdog Web Panel
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/gcp_watchdog
ExecStart=/opt/gcp_watchdog/venv/bin/gunicorn -w 1 -b 127.0.0.1:8765 app:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
